# start.ini
[program:EarthquakeSpy_12998]
autorestart=True
autostart=True
redirect_stderr=True
command=python3 manage.py runserver 0:12998 --insecure
user=root
directory=/root/applications/EarthquakeSpy
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /root/applications/logs/EarthquakeSpy.log
